# recipes-core/app-services/files/application.service
[Unit]
Description=Application service
#DefaultDependencies=no

[Service]
ExecStart=/usr/bin/application.sh

[Install]
WantedBy=multi-user.target
