# recipes-core/udev/files/10-console.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="g_serial",                         
RUN+="sleep 10 && systemctl enable getty@ttyGS0.service"

