# recipes-core/udev/files/50-fbcp.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="fb_ili9486", 
RUN+="/usr/bin/fbcp"

