# recipes-core/udev/files/50-fbcp.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="ads7846", 
RUN+="systemctl enable fbcp && systemctl start fbcp"

