# recipes-core/udev/files/55-psplash.rules
ACTION=="add", SUBSYSTEM=="module", KERNEL=="fb_ili9486", 
RUN+="psplash-start.service"


