# recipes-core/services/files/application.service
[Unit]
Description=Application service
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/application

[Install]
WantedBy=basic.target
